# scripts/discipline.service
[Unit]
Description=Computer time control service for kids

[Service]
ExecStart=/usr/local/bin/discipline

[Install]
WantedBy=multi-user.target